# deploy/shelter.conf
# Ubuntu upstart file at /etc/init/gunicorn.conf

description "Shelter System"
start on runlevel [2345]
stop on runlevel [!2345]

respawn
respawn limit 10 5
exec /srv/Shelter/deploy/shelter.sh
